# supervisord.prod.ini
[supervisord]
logfile=logs/supervisord.log
pidfile=supervisord.pid
logfile_maxbytes=50MB
logfile_backups=10
loglevel=debug
nodaemon=true
silent=false
minfds=1024
minprocs=200
nocleanup=false
strip_ansi=true
environment=PYTHONUNBUFFERED="1"

[unix_http_server]
file=supervisor.sock
chmod=0700
chown=%(ENV_USER)s:%(ENV_USER)s

[inet_http_server]
port=*:10000

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:virtmic]
command=bash -c "if [ ! -z ${STUDX_VIRTMIC+x} ]; then ./play.sh audio.mp3; fi"
process_name=%(program_name)s
directory=scripts
numprocs=1
priority=1
autostart=false
startsecs=1
startretries=5
autorestart=false
stdout_logfile=logs/virtmic-stdout.log
stderr_logfile=logs/virtmic-stderr.log

[program:writepad]
command=pnpm serve
process_name=%(program_name)s
directory=writepad
numprocs=1
priority=2
autostart=true
startsecs=1
startretries=5
autorestart=true
stdout_logfile=logs/writepad-stdout.log
stderr_logfile=logs/writepad-stderr.log

[program:whiteboard]
command=pnpm start:dev
process_name=%(program_name)s
directory=whiteboard
numprocs=1
priority=3
autostart=true
startsecs=1
startretries=5
autorestart=true
stdout_logfile=logs/whiteboard-stdout.log
stderr_logfile=logs/whiteboard-stderr.log

[program:ui]
command=pnpm dev
process_name=%(program_name)s
directory=ui
numprocs=1
priority=4
autostart=true
startsecs=1
startretries=5
autorestart=true
stdout_logfile=logs/ui-stdout.log
stderr_logfile=logs/ui-stderr.log

[program:api]
command=bash -c "source ../.env && poetry run -- ./manage.py runserver 0.0.0.0:$STUDX_API_PORT"  
process_name=%(program_name)s
directory=api
numprocs=1
priority=5
autostart=true
startsecs=1
startretries=5
autorestart=true
stdout_logfile=logs/api-stdout.log
stderr_logfile=logs/api-stderr.log

[program:celery]
command=bash -c "source ../.env && poetry run -- celery -A core worker"  
process_name=%(program_name)s
directory=api
numprocs=1
priority=6
autostart=true
startsecs=1
startretries=5
autorestart=true
stdout_logfile=logs/celery-stdout.log
stderr_logfile=logs/celery-stderr.log

[program:rooms]
command=bash -c "cargo watch -x 'run'"
process_name=%(program_name)s
directory=rooms
numprocs=1
priority=7
autostart=true
startsecs=1
startretries=5
autorestart=true
stdout_logfile=logs/rooms-stdout.log
stderr_logfile=logs/rooms-stderr.log

[program:web]
command=bash -c "if [[ ! -f Caddyfile ]]; then make caddyfile; fi; source .env && caddy run --watch"
process_name=%(program_name)s
numprocs=1
priority=8
autostart=true
startsecs=1
startretries=5
autorestart=true
stdout_logfile=logs/web-stdout.log
stderr_logfile=logs/web-stderr.log
